# delivery/systemd/carla-backend.service
[Unit]
Description=CARLA Delivery Backend
After=network-online.target carla-engine.service
Wants=network-online.target
Requires=carla-engine.service

[Service]
Type=simple
WorkingDirectory=%h/Carla_Project
ExecStart=/usr/bin/env bash %h/Carla_Project/delivery/start_backend.sh
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=default.target
